# cluster/sv/kubernetes-scheduler/run
#!/bin/sh

exec 2>&1

. /srv/kubernetes/etc/base.conf

exec "${BASE}/bin/kube-scheduler" --logtostderr --kubeconfig="${KUBECFG}"
