# cluster/sv/kubernetes-proxy/run
#!/bin/sh

exec 2>&1

. /srv/kubernetes/etc/base.conf

exec "${BASE}/bin/kube-proxy" --logtostderr --kubeconfig="${KUBECFG}" --master="https://144.76.60.131:6443"
